# etc_udev_rules.d/99-Gnuker.rules
KERNEL=="sd*[!0-9]", SUBSYSTEM=="block", ACTION=="add", RUN+="/usr/bin/sudo /usr/local/bin/Gnuke.sh $name"
